# deploy/systemd/happy-news@evening.timer
[Unit]
Description=Happy News Publisher - Evening (18:00)

[Timer]
OnCalendar=*-*-* 18:00:00
AccuracySec=1min
Persistent=false

[Install]
WantedBy=timers.target
