# deploy/systemd/happy-news@morning.timer
[Unit]
Description=Happy News Publisher - Morning (09:00)

[Timer]
OnCalendar=*-*-* 09:00:00
AccuracySec=1min
Persistent=false

[Install]
WantedBy=timers.target
