# deploy/systemd/happy-news@afternoon.timer
[Unit]
Description=Happy News Publisher - Afternoon (13:00)

[Timer]
OnCalendar=*-*-* 13:00:00
AccuracySec=1min
Persistent=false

[Install]
WantedBy=timers.target
